# bluealsa/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/dbus-1/
drwxr-xr-x	root/root	etc/dbus-1/system.d/
-rw-r--r--	root/root	etc/dbus-1/system.d/bluealsa.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/bluealsa
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/bluealsa
-rwxr-xr-x	root/root	usr/bin/bluealsa-aplay
drwxr-xr-x	root/root	usr/etc/
drwxr-xr-x	root/root	usr/etc/alsa/
drwxr-xr-x	root/root	usr/etc/alsa/conf.d/
-rw-r--r--	root/root	usr/etc/alsa/conf.d/20-bluealsa.conf
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/alsa-lib/
-rw-r--r--	root/root	usr/lib/alsa-lib/libasound_module_ctl_bluealsa.a
-rwxr-xr-x	root/root	usr/lib/alsa-lib/libasound_module_ctl_bluealsa.la
-rwxr-xr-x	root/root	usr/lib/alsa-lib/libasound_module_ctl_bluealsa.so
-rw-r--r--	root/root	usr/lib/alsa-lib/libasound_module_pcm_bluealsa.a
-rwxr-xr-x	root/root	usr/lib/alsa-lib/libasound_module_pcm_bluealsa.la
-rwxr-xr-x	root/root	usr/lib/alsa-lib/libasound_module_pcm_bluealsa.so
